# templates/airflow-flower.service.erb
[Unit]
Description= Airflow flower
Wants=basic.target
After=basic.target network.target

[Service]
Environment="AIRFLOW_HOME=<%= scope.lookupvar('airflow::home_folder')%>"
User=<%= scope.lookupvar('airflow::user')%>
Group=<%= scope.lookupvar('airflow::group')%>
Type=simple
ExecStart=/bin/airflow flower
KillMode=process
Restart=on-failure
RestartSec=15s

[Install]
WantedBy=multi-user.target
